# etc/profile.d/motd.sh
#!/bin/bash
echo -e "\e[38;5;45m$( < /etc/motd.logo )\e[0m\n"

