# etc/udev/rules.d/usbdac.rules
# udevadm control --reload-rules && udevadm trigger
ACTION=="add", KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/srv/http/settings/mpd-conf.sh"
ACTION=="remove", KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/srv/http/settings/mpd-conf.sh remove"