# etc/udev/rules.d/usbdac.rules
ACTION=="add|remove", KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/srv/http/settings/mpdconf.sh"
