# etc/udev/rules.d/usbdac.rules
ACTION=="add", KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/srv/http/settings/mpdconf.sh"
ACTION=="remove", KERNEL=="card*", SUBSYSTEM=="sound", RUN+="/srv/http/settings/mpdconf.sh remove"