# etc/systemd/system/nginx.service.d/privatetmp-false.conf
# allow mount for all users
[Service]
PrivateTmp=false
