# etc/systemd/system/nginx.service.d/privatetmp-false.conf
[Unit]
After=php-fpm.service

# allow mount for all users
[Service]
PrivateTmp=false
