# etc/systemd/system/libraryrandom.service
[Unit]
Description=Library random
After=mpd.service

[Service]
ExecStart=/srv/http/libraryrandom.sh
