# etc/systemd/system/mpdidle.service
[Unit]
Description=MPD idle
After=mpd.service

[Service]
ExecStart=/srv/http/mpdidle.sh
