# etc/systemd/system/radiowatchdog.service
[Unit]
Description=MPD watchdog

[Service]
ExecStart=/srv/http/radiowatchdog.sh
