# etc/systemd/system/runonce.service
[Unit]
Description=RuneUIe run once
ConditionDirectoryNotEmpty=/mnt/MPD/USB

[Service]
Type=oneshot
ExecStart=/srv/http/runonce.sh

[Install]
WantedBy=multi-user.target
