# etc/systemd/system/startup.service
[Unit]
Description=RuneUIe Startup
After=php-fpm.service

[Service]
Type=oneshot
ExecStart=/srv/http/startup.sh

[Install]
WantedBy=multi-user.target
