# etc/systemd/system/usbdacoff.service
[Unit]
Description=Hotplug USB DAC
After=multi-user.target
[Service]
Type=oneshot
ExecStart=/usr/local/bin/usbdac.sh
