# etc/systemd/system/usbdac.service
[Unit]
Description=Hotplug USB DAC

[Service]
Type=oneshot
ExecStart=/srv/http/settings/mpdconf.sh
