# etc/systemd/system/bluetooth.service.d/execstartpost-discoverable.conf
[Service]
ExecStartPost=/srv/http/settings/systembluetooth.sh
