# etc/systemd/system/bluetooth.service.d/override.conf
[Unit]
After=startup.service

[Service]
ExecStartPost=/srv/http/settings/system-bluetooth.sh

